# roles/docker_host/templates/docker-prune.service.j2
[Unit]
Description=Prune unused Docker data
Wants=docker.service
After=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/docker-prune
